# tools/systemd/noctis-pro-tunnel.service
[Unit]
Description=Noctis Pro Public HTTPS Tunnel (ngrok)
After=network-online.target noctis-pro.service
Wants=network-online.target

[Service]
Type=simple
Environment=APP_DIR=/opt/noctispro
WorkingDirectory=%E{APP_DIR}
# Optional env files used by server deploy scripts
EnvironmentFile=-/etc/noctis-pro/noctis-pro.env
EnvironmentFile=-/etc/noctispro/noctispro.env

# Reserved domain you want to access (override in env file if needed)
Environment=NGROK_DOMAIN=noctis-pro.com.ngrok.app

# ngrok config path (written server-side; MUST NOT be committed)
Environment=NGROK_CONFIG=/etc/noctis-pro/ngrok.yml
Environment=NGROK_BIN=/usr/local/bin/ngrok

# Persist tunnel URL so Django can auto-detect it (see noctis_pro/settings.py)
ExecStartPre=/usr/bin/env bash -lc 'if [[ -n "${NGROK_DOMAIN:-}" ]]; then echo -n "https://${NGROK_DOMAIN}" > "%E{APP_DIR}/.tunnel-url"; fi'

# Start ngrok (fallback to /etc/ngrok.yml and PATH if needed)
ExecStart=/usr/bin/env bash -lc 'BIN="${NGROK_BIN:-/usr/local/bin/ngrok}"; if [[ ! -x "${BIN}" ]]; then BIN="$(command -v ngrok || true)"; fi; if [[ -z "${BIN}" ]]; then echo "ERROR: ngrok binary not found (install ngrok to /usr/local/bin/ngrok or PATH)" >&2; exit 127; fi; CFG="${NGROK_CONFIG:-/etc/noctis-pro/ngrok.yml}"; if [[ ! -f "${CFG}" && -f "/etc/ngrok.yml" ]]; then CFG="/etc/ngrok.yml"; fi; exec "${BIN}" start --config "${CFG}" noctis-web'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
