# tools/systemd/noctis-pro-dicom.service
[Unit]
Description=Noctis Pro PACS DICOM Receiver (SCP)
After=network-online.target

[Service]
Type=simple
User=noctispro
Group=noctispro
EnvironmentFile=-/etc/noctis-pro/noctis-pro.env
EnvironmentFile=-/etc/noctispro/noctispro.env
WorkingDirectory=/
ExecStart=/usr/bin/env bash -lc 'APP_DIR="${APP_DIR:-/opt/noctispro}"; cd "${APP_DIR}"; exec "${APP_DIR}/venv/bin/python" "${APP_DIR}/dicom_receiver.py" --port 11112 --aet NOCTIS_SCP'
Restart=always
RestartSec=5

# --- systemd hardening (baseline) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target