# tools/systemd/noctis-pro-tunnel.service
[Unit]
Description=Noctis Pro Public HTTPS Tunnel (ngrok)
After=network-online.target noctis-pro.service
Wants=network-online.target

[Service]
Type=simple
User=noctispro
Group=noctispro
Environment=APP_DIR=/opt/noctispro
# Avoid relying on newer systemd specifiers for env expansion here; we will `cd`
# inside the shell commands so this unit works across more distros.
WorkingDirectory=/
# Optional env files used by server deploy scripts
EnvironmentFile=-/etc/noctis-pro/noctis-pro.env
EnvironmentFile=-/etc/noctispro/noctispro.env

# Optional reserved domain to use for the tunnel.
# IMPORTANT: do not hardcode a domain here unless you have actually reserved it
# in your ngrok account; otherwise the tunnel may fail to start.
# Set this in /etc/noctis-pro/noctis-pro.env if you have one:
#   NGROK_DOMAIN=your-reserved-domain.ngrok.app

# ngrok config path (written server-side; MUST NOT be committed)
Environment=NGROK_CONFIG=/etc/noctis-pro/ngrok.yml
Environment=NGROK_BIN=/usr/local/bin/ngrok

# Require a pre-created ngrok config file (deploy.sh writes it with safe perms).
ExecStartPre=/usr/bin/env bash -lc 'CFG="$${NGROK_CONFIG:-/etc/noctis-pro/ngrok.yml}"; if [[ ! -f "$${CFG}" ]]; then echo "ERROR: Missing $${CFG}. Create it (or re-run deploy.sh in ngrok mode)." >&2; exit 1; fi'

# Persist tunnel URL so Django can auto-detect it (see noctis_pro/settings.py)
ExecStartPre=/usr/bin/env bash -lc 'cd "$${APP_DIR:-/opt/noctispro}" 2>/dev/null || exit 0; if [[ -n "$${NGROK_DOMAIN:-}" ]]; then printf "%s" "https://$${NGROK_DOMAIN}" > ".tunnel-url" 2>/dev/null || true; fi; exit 0'

# Start ngrok (fallback to /etc/ngrok.yml and PATH if needed)
ExecStart=/usr/bin/env bash -lc 'BIN="$${NGROK_BIN:-/usr/local/bin/ngrok}"; if [[ ! -x "$${BIN}" ]]; then BIN="$$(command -v ngrok 2>/dev/null || true)"; fi; if [[ -z "$${BIN}" ]]; then echo "ERROR: ngrok binary not found (install ngrok to /usr/local/bin/ngrok or PATH)" >&2; exit 127; fi; CFG="$${NGROK_CONFIG:-/etc/noctis-pro/ngrok.yml}"; if [[ ! -f "$${CFG}" && -f "/etc/ngrok.yml" ]]; then CFG="/etc/ngrok.yml"; fi; exec "$${BIN}" start --config "$${CFG}" noctis-web'

# If no reserved domain is provided, capture the randomly assigned public URL
# from the local ngrok API and persist it for Django.
ExecStartPost=/usr/bin/env bash -lc 'cd "$${APP_DIR:-/opt/noctispro}" 2>/dev/null || exit 0; if [[ -z "$${NGROK_DOMAIN:-}" ]]; then if ! command -v curl >/dev/null 2>&1 || ! command -v python3 >/dev/null 2>&1; then exit 0; fi; for _ in $$(seq 1 60); do URL="$$(curl -fsS http://127.0.0.1:4040/api/tunnels 2>/dev/null | python3 - <<'"'"'PY'"'"'\nimport json,sys\ntry:\n    data=json.load(sys.stdin)\nexcept Exception:\n    sys.exit(0)\nfor t in data.get(\"tunnels\", []) or []:\n    u=(t.get(\"public_url\") or \"\").strip()\n    if u.startswith(\"https://\"):\n        print(u)\n        break\nPY\n)"; if [[ -n "$${URL}" ]]; then printf "%s" "$${URL}" > ".tunnel-url" 2>/dev/null || true; break; fi; sleep 0.5; done; fi; exit 0'
Restart=always
RestartSec=5

# --- systemd hardening (baseline) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
