# tools/systemd/noctis-pro-tunnel.service
[Unit]
Description=Noctis Pro Public HTTPS Tunnel (ngrok)
After=network-online.target noctis-pro.service
Wants=network-online.target

[Service]
Type=simple
Environment=APP_DIR=/opt/noctispro
WorkingDirectory=%E{APP_DIR}
# Optional env files used by server deploy scripts
EnvironmentFile=-/etc/noctis-pro/noctis-pro.env
EnvironmentFile=-/etc/noctispro/noctispro.env

# Optional reserved domain to use for the tunnel.
# IMPORTANT: do not hardcode a domain here unless you have actually reserved it
# in your ngrok account; otherwise the tunnel may fail to start.
# Set this in /etc/noctis-pro/noctis-pro.env if you have one:
#   NGROK_DOMAIN=your-reserved-domain.ngrok.app

# ngrok config path (written server-side; MUST NOT be committed)
Environment=NGROK_CONFIG=/etc/noctis-pro/ngrok.yml
Environment=NGROK_BIN=/usr/local/bin/ngrok

# If NGROK_AUTHTOKEN is set (in /etc/noctis-pro/noctis-pro.env), ensure the ngrok
# config exists. This keeps the token out of git while making setup one-step.
ExecStartPre=/usr/bin/env bash -lc 'CFG="${NGROK_CONFIG:-/etc/noctis-pro/ngrok.yml}"; if [[ ! -f "${CFG}" ]]; then if [[ -z "${NGROK_AUTHTOKEN:-}" ]]; then echo "ERROR: ${CFG} missing and NGROK_AUTHTOKEN not set. Put NGROK_AUTHTOKEN in /etc/noctis-pro/noctis-pro.env or create ${CFG}." >&2; exit 1; fi; install -d -m 0755 "$(dirname "${CFG}")"; umask 077; { echo "version: 2"; echo "authtoken: ${NGROK_AUTHTOKEN}"; echo "tunnels:"; echo "  noctis-web:"; echo "    proto: http"; echo "    addr: 8000"; echo "    schemes:"; echo "      - https"; if [[ -n "${NGROK_DOMAIN:-}" ]]; then echo "    domain: ${NGROK_DOMAIN}"; fi; } > "${CFG}"; chmod 600 "${CFG}"; fi'

# Persist tunnel URL so Django can auto-detect it (see noctis_pro/settings.py)
ExecStartPre=/usr/bin/env bash -lc 'if [[ -n "${NGROK_DOMAIN:-}" ]]; then echo -n "https://${NGROK_DOMAIN}" > "%E{APP_DIR}/.tunnel-url"; fi'

# Start ngrok (fallback to /etc/ngrok.yml and PATH if needed)
ExecStart=/usr/bin/env bash -lc 'BIN="${NGROK_BIN:-/usr/local/bin/ngrok}"; if [[ ! -x "${BIN}" ]]; then BIN="$(command -v ngrok || true)"; fi; if [[ -z "${BIN}" ]]; then echo "ERROR: ngrok binary not found (install ngrok to /usr/local/bin/ngrok or PATH)" >&2; exit 127; fi; CFG="${NGROK_CONFIG:-/etc/noctis-pro/ngrok.yml}"; if [[ ! -f "${CFG}" && -f "/etc/ngrok.yml" ]]; then CFG="/etc/ngrok.yml"; fi; exec "${BIN}" start --config "${CFG}" noctis-web'

# If no reserved domain is provided, capture the randomly assigned public URL
# from the local ngrok API and persist it for Django.
ExecStartPost=/usr/bin/env bash -lc 'URL_FILE="%E{APP_DIR}/.tunnel-url"; if [[ -z "${NGROK_DOMAIN:-}" ]]; then for _ in $(seq 1 60); do URL="$(curl -fsS http://127.0.0.1:4040/api/tunnels 2>/dev/null | python3 - <<'"'"'PY'"'"'\nimport json,sys\ntry:\n    data=json.load(sys.stdin)\nexcept Exception:\n    sys.exit(0)\nfor t in data.get(\"tunnels\", []) or []:\n    u=(t.get(\"public_url\") or \"\").strip()\n    if u.startswith(\"https://\"):\n        print(u)\n        break\nPY\n)"; if [[ -n "${URL}" ]]; then echo -n "${URL}" > "${URL_FILE}"; break; fi; sleep 0.5; done; fi'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
